# com.wirevpn.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.wirevpn.startup.plist — launchd daemon for WireGuard persistence
  Installs to: /Library/LaunchDaemons/com.wirevpn.startup.plist

  This ensures the VPN tunnel auto-connects on every Mac boot.
  Without this, you would need to manually run wg-quick up after each restart.

  Logs: /var/log/wirevpn.log
  To install: sudo launchctl load /Library/LaunchDaemons/com.wirevpn.startup.plist
  To remove:  sudo launchctl unload /Library/LaunchDaemons/com.wirevpn.startup.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wirevpn.startup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/wg-quick</string>
        <string>up</string>
        <string>/Users/Shared/WireVPN/client.conf</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/var/log/wirevpn.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/wirevpn.log</string>
</dict>
</plist>
